# etc/cron.d/rhizomatica
#Rhizomatica cron jobs
PYTHONPATH=$PYTHONPATH:/var/rhizomatica/rccn

# m h  dom mon dow   user command
  0 8  1   *   *     root /usr/bin/python /var/rhizomatica/rccn/rsc.py notice
  0 8  3   *   *     root /usr/bin/python /var/rhizomatica/rccn/rsc.py reminder
  0 21 7   *   *     root /usr/bin/python /var/rhizomatica/rccn/rsc.py deactivate

